# c2servers/cron.d/redelk_cobaltstrike
#
# Part of RedELK
# cron.d script for periodic actions related to RedELK and Cobalt Strike
#
# Author: Outflank B.V. / Marc Smeets 
#

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Command to sync the logs from cobaltstrike to our scponly user's home directory
# m h dom mon dow user  command
* * * * * root /usr/bin/rsync -ax --delete /root/cobaltstrike/logs /home/scponly/; /bin/chown -R scponly:scponly /home/scponly/*
* * * * * root /usr/bin/rsync -ax --delete /root/cobaltstrike/profiles /home/scponly/; /bin/chown -R scponly:scponly /home/scponly/*
* * * * * root /usr/bin/rsync -ax --delete /root/cobaltstrike/data /home/scponly/; /bin/chown -R scponly:scponly /home/scponly/*
* * * * * root /usr/share/redelk/bin/export_cobaltstrikedata.sh
* * * * * root /usr/share/redelk/bin/copydownloads_cobaltstrike.sh
